# scripts/capture.cron
# this file goes in /etc/cron.d

*/5 * * * *	pi	/usr/local/bin/capture.sh >> /var/log/pi-env/captures.log 2>&1
